# server/deploy/stockdata-worker.service
[Unit]
Description=stockdata worker（Celery，solo pool，串行执行抓取任务）
Requires=postgresql.service valkey.service
After=network.target postgresql.service valkey.service

[Service]
Type=exec
User=david
WorkingDirectory=/home/david/workspace/stockdata/server
ExecStart=/usr/bin/env uv run celery -A fetcher.app worker --pool=solo --concurrency=1 --loglevel=info
TimeoutStopSec=120
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
