# deploy/stockdata-mcp.service
[Unit]
Description=stockdata MCP shim (FastMCP → REST 转发)
Requires=stockdata-api.service
After=network.target stockdata-api.service

[Service]
Type=exec
User=david
WorkingDirectory=/home/david/workspace/stockdata
ExecStart=/usr/bin/env uv run --group mcp python -m mcp_shim.main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
